# restic/systemd/restic-all.service

[Unit]
Description=Restic backup everything service

[Service]
Type=oneshot
ExecStart=/home/jdheyburn/dotfiles/restic/restic-all.sh
