# restic/systemd/restic-prune.service
[Unit]
Description=Restic backup service (data pruning)

[Service]
Type=oneshot
ExecStart=/home/jdheyburn/dotfiles/restic/restic-all.sh prune
